# config/supervisor.conf
[program:myproject]
command=/home/ubuntu/ChatApp/venv/bin/daphne -p 8001 config.asgi:application
directory=/home/ubuntu/ChatApp/
user=ubuntu
autorestart=true
redirect_stderr=true
